# packages/mcp/systemd/litopys-mcp.service
[Unit]
Description=Litopys MCP Server (HTTP/SSE)
Documentation=https://github.com/litopys-dev/litopys
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/litopys
Environment=LITOPYS_MCP_PORT=7777
Environment=LITOPYS_GRAPH_PATH=%h/.litopys/graph
EnvironmentFile=%h/.litopys/mcp.env
ExecStart=%h/.bun/bin/bun packages/mcp/src/http.ts
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
