# packages/mcp/systemd/litopys-mcp.service
[Unit]
Description=Litopys MCP Server (HTTP/SSE)
After=network.target

[Service]
Type=simple
User=denis
WorkingDirectory=/home/denis/litopys
Environment=LITOPYS_MCP_PORT=7777
Environment=LITOPYS_GRAPH_PATH=/home/denis/.litopys/graph
EnvironmentFile=/home/denis/.litopys/mcp.env
ExecStart=/home/denis/.bun/bin/bun packages/mcp/src/http.ts
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
